# polybar/.config/systemd/user/check-updates.timer
[Unit]
Description=Count how many updates are available

[Timer]
OnCalendar=hourly
Unit=check-updates.service

[Install]
WantedBy=timers.target
